# awg-client-package/etc/xdg/autostart/awg-client.desktop
[Desktop Entry]
Type=Application
Name=AWGuird AutoRun
Comment=Start AWGuird VPN interface monitor on session login
Exec=pkexec /usr/local/bin/awg-client
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
